# Services/docker_podcast.service
[Unit]
Description=docker_podcast

[Service]
Type=oneshot
