# Services/docker_podcast.timer
[Unit]
Description=docker_podcast

[Timer]
OnCalendar=00/1:05
Persistent=true

[Install]
WantedBy=timers.target