# Services/docker_video.service
[Unit]
Description=docker_video

[Service]
Type=oneshot
